# roles/jukebox/templates/quadlets/jukebox-player.container.j2
[Unit]
Description=jukebox player (squeezelite)
PartOf=jukebox-pod.service

[Container]
ContainerName=%N
Image={{ jukebox_player_image }}
AutoUpdate=registry

Pod=jukebox.pod

# ALSA access — AddDevice passes the device node into the container,
# but SELinux blocks the rootless process from opening it.
# SecurityLabelDisable=true is the minimal override (equivalent to
# --security-opt label=disable). Without this, squeezelite sees the
# device nodes but gets "Permission denied" when opening them.
SecurityLabelDisable=true
AddDevice=/dev/snd
{% if asound_conf_stat is defined and asound_conf_stat.stat.exists %}
# Host has an /etc/asound.conf override (e.g. VM with playback on
# card != 0). Mount it so squeezelite picks the correct default card.
Volume=/etc/asound.conf:/etc/asound.conf:ro
{% endif %}
GroupAdd=keep-groups

Environment=PRESET={{ jukebox_squeezelite_preset }}
Environment=SQUEEZELITE_NAME={{ jukebox_squeezelite_name }}
{% if jukebox_squeezelite_mac %}
Environment=SQUEEZELITE_MAC_ADDRESS={{ jukebox_squeezelite_mac }}
{% endif %}
#Environment=STARTUP_DELAY_SEC=10000
#Environment=SQUEEZELITE_LOG_CATEGORY_ALL=debug
