# roles/jukebox/templates/quadlets/jukebox-player.container.j2
[Unit]
Description=jukebox player (squeezelite)
PartOf=jukebox-pod.service

[Container]
ContainerName=%N
Image=docker.io/giof71/squeezelite
AutoUpdate=registry

Pod=jukebox.pod

# ALSA access (adjust if needed)
AddDevice=/dev/snd
{% if asound_conf_stat is defined and asound_conf_stat.stat.exists %}
# Host has an /etc/asound.conf override (e.g. VM with playback on
# card != 0). Mount it so squeezelite picks the correct default card —
# without this, it crashes with "playback open error: No such file or
# directory" and the whole pod is torn down.
Volume=/etc/asound.conf:/etc/asound.conf:ro
{% endif %}
GroupAdd=keep-groups

Environment=PRESET=topping-D10s
#Environment=STARTUP_DELAY_SEC=10000
#Environment=SQUEEZELITE_LOG_CATEGORY_ALL=debug
