# roles/backup/files/home-server-backup.service
[Unit]
Description=Home server volume backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/home-server-backup.sh
StandardOutput=journal
StandardError=journal
