# roles/paperless-ngx/templates/quadlets/paperless-ngx.container.j2
[Unit]
Description=paperless-ngx web application
PartOf=paperless-ngx-pod.service
After=paperless-db.service
After=paperless-redis.service

[Container]
ContainerName=%N
Image={{ paperless_app_image }}
AutoUpdate=registry

Pod=paperless-ngx.pod

Volume=paperless-media.volume:/usr/src/paperless/media
Volume=paperless-data.volume:/usr/src/paperless/data
Volume=paperless-export.volume:/usr/src/paperless/export
Volume=paperless-consume.volume:/usr/src/paperless/consume

# Database
Environment=PAPERLESS_DBHOST=localhost
Environment=PAPERLESS_DBPORT=5432
Environment=PAPERLESS_DBNAME=paperless
Environment=PAPERLESS_DBUSER=paperless
Environment=PAPERLESS_DBPASS={{ paperless_db_password }}

# Redis
Environment=PAPERLESS_REDIS=redis://localhost:6379

# Security
Environment=PAPERLESS_SECRET_KEY={{ paperless_secret_key }}

# Admin account (created on first startup only)
Environment=PAPERLESS_ADMIN_USER={{ paperless_admin_user }}
Environment=PAPERLESS_ADMIN_PASSWORD={{ paperless_admin_password }}

# Localisation
Environment=PAPERLESS_TIME_ZONE={{ paperless_time_zone }}
Environment=PAPERLESS_OCR_LANGUAGE={{ paperless_ocr_language }}

# Gotenberg (document conversion)
{% if paperless_enable_gotenberg %}
Environment=PAPERLESS_TIKA_ENABLED=true
Environment=PAPERLESS_TIKA_GOTENBERG_ENDPOINT=http://localhost:3000
Environment=PAPERLESS_TIKA_ENDPOINT=http://localhost:3000
{% endif %}

# Misc
Environment=PAPERLESS_URL={{ paperless_url | default('') }}
Environment=USERMAP_UID=1000
Environment=USERMAP_GID=1000

# Consumer behaviour
Environment=PAPERLESS_CONSUMER_RECURSIVE={{ paperless_consumer_recursive | bool | string | lower }}
Environment=PAPERLESS_CONSUMER_SUBDIRS_AS_TAGS={{ paperless_consumer_subdirs_as_tags | bool | string | lower }}
Environment=PAPERLESS_CONSUMER_DELETE_DUPLICATES={{ paperless_consumer_delete_duplicates | bool | string | lower }}

HealthCmd=curl -fs http://localhost:8000 || exit 1
HealthInterval=30s
HealthStartPeriod=120s
HealthRetries=5

[Service]
TimeoutStartSec=600
