# roles/dashboard/files/home-server-webserver.service
[Unit]
Description=Home server dashboard web server
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/var/www/dashboard
ExecStart=/usr/bin/python3 -m http.server 8888 --bind 0.0.0.0
Restart=on-failure

[Install]
WantedBy=multi-user.target
