# roles/pihole/templates/quadlets/unbound.container.j2
[Unit]
Description=Unbound recursive DNS resolver (Pi-hole upstream)
After=network-online.target
Wants=network-online.target

[Container]
ContainerName=%N
Image=docker.io/klutchell/unbound:latest
AutoUpdate=registry
Network=host

# Mount our config and override the entrypoint to use it.
# The klutchell image's default entrypoint runs unbound with its
# own config (which binds port 53). We bypass it and run unbound
# directly with our config that only binds port {{ pihole_unbound_port }}.
Volume=/home/pihole/unbound-pihole.conf:/etc/unbound/pihole.conf:ro,Z
Entrypoint=unbound
Exec=-d -c /etc/unbound/pihole.conf

[Service]
Restart=on-failure
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target default.target
