# roles/pihole/templates/quadlets/unbound.container.j2
[Unit]
Description=Unbound recursive DNS resolver (Pi-hole upstream)
After=network-online.target pihole.service
Wants=network-online.target
Requires=pihole.service

[Container]
ContainerName=%N
Image=docker.io/klutchell/unbound:latest
AutoUpdate=registry
# Share Pi-hole's network namespace so both containers see the same
# loopback — Pi-hole can reach Unbound at 127.0.0.1:{{ pihole_unbound_port }}.
Network=container:pihole

# SELinux blocks rootless bind-mount reads without this.
# Same fix as jukebox/squeezelite. Container is still UID-confined.
SecurityLabelDisable=true

# Mount our config and bypass the image's entrypoint (which binds
# port 53, conflicting with Pi-hole). We run unbound directly with
# our config that only listens on port {{ pihole_unbound_port }}.
Volume=/home/pihole/unbound-pihole.conf:/etc/unbound/pihole.conf:ro
Entrypoint=unbound
Exec=-d -c /etc/unbound/pihole.conf

[Service]
Restart=on-failure
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target default.target
