# roles/backup/templates/home-server-backup.timer.j2
[Unit]
Description=Nightly home server backup

[Timer]
OnCalendar=*-*-* {{ backup_time }}
Persistent=true

[Install]
WantedBy=timers.target
